# machined-dns.service
[Unit]
Description=Machined dns service

[Service]
Type=simple
ExecStart=/usr/local/bin/machined-dns
